# remote-socket-sunset-control/crontab
*/10 * * * * php /home/pi/remote-socket/timer.php
